# templates/lab-backup.service.j2
[Unit]
Description=Lightweight backup for {{ inventory_hostname }}
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/lab-backup.sh
